# grafana_docker/supervisord/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[program:influxdb]
priority = 10
command = /usr/bin/influxd

[program:telegraf]
priority = 10
command = /usr/bin/telegraf

[program:grafana]
priority = 10
command = /usr/share/grafana/bin/grafana-server
autostart = true
autorestart = true
startretries = 3
user = grafana


